# services/DcpMqttClient/run
#!/bin/sh
exec 2>&1 
exec /data/DcpMqttClient/DcpMqttClient.py